# paperless_sync.service
[Unit]
Description=Sync Paperless archive to mounted SMB share
After=network.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/paperless_sync.sh
User=root
Group=root
Nice=10
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

# Optional: stdout/stderr an journal weiterleiten
StandardOutput=journal
StandardError=journal
